# modules/sr_site/templates/willie.service.erb
# Copied from willie's contrib, with slight changes
[Unit]
Description=Willie IRC bot
Documentation=http://willie.dftba.net/
After=network.target

[Service]
Type=simple
User=willie
# Run ExecStartPre with root-permissions so we can make the 'irc' folder
# which we use for our pipe
PermissionsStartOnly=true
ExecStartPre=-/usr/bin/mkdir <%= @pipe_dir %>
ExecStartPre=/usr/bin/chown -R willie <%= @pipe_dir %>
# Run ExecStart with User=willie
ExecStart=/usr/bin/willie -c <%= @config_path %> --quiet
ExecStop=/bin/kill -USR1 $MAINPID
PIDFile=/run/willie/willie.pid
Restart=on-failure
RestartPreventExitStatus=2
RestartSec=30
PrivateTmp=False
